# cron/bin/crontab
#!/bin/sh

env EDITOR=nano /data/hack/bin/busybox crontab -e

